# deploy/salt/roots/salt/twisted_gateway/supervisord.conf
[program:twisted_gateway]
command = /var/www/joltem_venv/bin/twistd --nodaemon -y gateway/gateway.tac
directory = /joltem/
user = root
autostart = true
autorestart = true
redirect_stderr = True
