# deploy/debian/supervisor.ini
[program:joltem-web]
autorestart=true
autostart=true
command=newrelic-admin run-program /usr/lib/joltem/venv/bin/uwsgi --ini /usr/lib/joltem/uwsgi.ini
directory=/usr/lib/joltem/joltem
environment=NEW_RELIC_CONFIG_FILE="/usr/lib/joltem/configuration/newrelic.ini"
redirect_stderr=true
stdout_logfile=/usr/lib/joltem/log/web.log
user=joltem

[program:joltem-worker]
autorestart=true
autostart=true
command=newrelic-admin run-program /usr/lib/joltem/venv/bin/celery -A joltem worker -l info
directory=/usr/lib/joltem/joltem
environment=NEW_RELIC_CONFIG_FILE="/usr/lib/joltem/configuration/newrelic.ini"
redirect_stderr=true
stdout_logfile=/usr/lib/joltem/log/worker.log
user=joltem

[program:joltem-beat]
autorestart=true
autostart=true
command=newrelic-admin run-program /usr/lib/joltem/venv/bin/celery -A joltem beat -l info --pidfile=/usr/lib/joltem/run/beat.pid -s /usr/lib/joltem/run/beat-schedule
directory=/usr/lib/joltem/joltem
environment=NEW_RELIC_CONFIG_FILE="/usr/lib/joltem/configuration/newrelic.ini"
redirect_stderr=true
stdout_logfile=/usr/lib/joltem/log/beat.log
user=joltem

[program:joltem-gateway]
autorestart=true
autostart=true
command=/usr/lib/joltem/venv/bin/twistd --nodaemon -y gateway/gateway.tac --pidfile=/usr/lib/joltem/run/twistd.pid
directory=/usr/lib/joltem/joltem
redirect_stderr=true
stdout_logfile=/usr/lib/joltem/log/gateway.log
